# deploy/serving/docker/assets/configurator_dumpenv.service
[Unit]
Description=Dump the Docker environment variables

[Service]
Type=oneshot
ExecStart=/root/.system/configurator_dumpenv.sh
TimeoutSec=0
